# systemd/user/emoji-picker-addon.service
[Unit]
Description=Reapply the Quickshell emoji picker addon

[Service]
Type=oneshot
ExecStart=%h/.local/share/quickshell-addons/emoji-picker/apply.sh
