# systemd/user/headset-mic-loopback-addon.service
[Unit]
Description=Reapply the headset audio virtual microphone control

[Service]
Type=oneshot
ExecStart=%h/.local/share/quickshell-addons/headset-mic-loopback/apply.sh
