# systemd/user/idle-inhibit-addon.path
[Unit]
Description=Watch the Hyprland idle setting

[Path]
PathChanged=%h/.config/hypr/settings.json
Unit=idle-inhibit-addon.service

[Install]
WantedBy=default.target
