# systemd/user/launcher-web-search-addon.path
[Unit]
Description=Watch the Quickshell application launcher extended with Zen search

[Path]
PathChanged=%h/.config/hypr/scripts/quickshell/applauncher/appLauncher.qml
Unit=launcher-web-search-addon.service

[Install]
WantedBy=default.target
