# systemd/user/matugen-vibrant-addon.service
[Unit]
Description=Reapply the vibrant Matugen colors addon

[Service]
Type=oneshot
ExecStart=%h/.local/share/quickshell-addons/matugen-vibrant/apply.sh
