# systemd/user/music-preview-rounded-addon.service
[Unit]
Description=Reapply rounded music preview artwork

[Service]
Type=oneshot
ExecStart=%h/.local/share/quickshell-addons/music-preview-rounded/apply.sh
