# systemd/user/screenshot-freeze-addon.service
[Unit]
Description=Reapply the frozen screenshot addon

[Service]
Type=oneshot
ExecStart=%h/.local/share/quickshell-addons/screenshot-freeze/apply.sh
