# systemd/user/wallpaper-random-addon.service
[Unit]
Description=Reapply the Quickshell random wallpaper addon

[Service]
Type=oneshot
ExecStart=%h/.local/share/quickshell-addons/wallpaper-random/apply.sh
