# systemd/user/zoomit-addon.service
[Unit]
Description=Reapply the Hyprland ZoomIt addon

[Service]
Type=oneshot
ExecStart=%h/.local/share/quickshell-addons/zoomit/apply.sh
